# recipebook.service
[Unit]
Description=uWSGI instance to serve RecipeBook
After=network.target

[Service]
User=
